# files.list
/lib/modules/<KVERSION>/misc/dgrp.ko
/usr/share/man/man8/dgrp.8
/usr/sbin/dgelreset
/usr/sbin/dgipserv
/usr/share/man/man1/ditty-rp.1
/usr/share/man/man1/dgipserv.1
/usr/share/man/man1/dgelreset.1
/usr/share/man/man8/dgrp_gui.8
/usr/share/man/man8/dgrp_cfg_node.8
/usr/share/man/man8/dgrp.8
/usr/share/man/man8/drpd.8
/usr/bin/ditty-rp
/usr/bin/dinc.dgrp
/usr/bin/dgrp/daemon/drpd
/usr/bin/dgrp/config/dgrp_gui
/usr/bin/dgrp/config/postinstall
/usr/bin/dgrp/config/file_locations
/usr/bin/dgrp/config/dgrp_cfg_node
/usr/bin/dgrp/config/dgrp_mk_spec
/usr/bin/dgrp/config/preuninstall
/usr/bin/dgrp/config/dgrp.gif
/usr/bin/dgrp/config/dgrp_mk_specs
/usr/bin/dgrp_udev
/usr/bin/dpa.dgrp
/etc/rc.d/init.d/dgrp_daemon
/etc/rc.d/init.d/dgrp_ditty


